# ansible/roles/k3s/prereq/templates/cilium-watchdog.service.j2
[Unit]
Description=Cilium migration dead-man switch (SSH connectivity watchdog)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cilium-watchdog.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
